# lib/cap_recipes/tasks/geoip/geoip.cron.erb
# 11:00am UTC ==> 3:00am PST on Wed
0 11 * * 3 root /usr/bin/geoipupdate
